# contrib/systemd/stamp-reconcile.timer
[Unit]
Description=Run stamp reconcile daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
