# contrib/systemd/stamp-reconcile.service
[Unit]
Description=Stamp package drift reconciliation

[Service]
ExecStart=%h/.local/bin/stamp reconcile
Type=oneshot
